# usb-display-switch.service
[Unit]
Description=usb-display-switch

[Service]
ExecStart=%h/.local/bin/usb-display-switch
Restart=always

[Install]
WantedBy=default.target
